# InfiniteApiRunner.service
[Unit]
Description=The Infinite Experiment bot

[Service]
User=root
WorkingDirectory=/root/projects/infinite-experiment-api
ExecStart=/bin/bash -c "npm run build && npm start:prod"
Restart=always

[Install]
WantedBy=multi-user.target
